# laptop-backup/systemd_units/exthdd.automount
[Unit]
Description="Automounting of USB External HDD"

[Automount]
Where=/exthdd/
TimeoutIdleSec=5min

[Install]
WantedBy=local-fs.target
RequiredBy=borg_backup.service