# laptop-backup/deployment/roles/laptop2hdd/templates/borg_backup.timer.j2
[Unit]
Description="Run laptop backup every week"

[Timer]
OnCalendar={{ backup_time }}
Persistent=true

[Install]
WantedBy=timers.target
